# testing/crontab

crontab -e
@reboot sudo python /home/pi/timelapse/raspiLapseCam.py &

crontab -l